# conf/balloon.cron
# minutes hours day1-31 month1-12 day1-7 command
0 3 * * * /home/balloon/backend/manage.py forecast_download > /home/balloon/log/download-$(date +%Y-%m-%dT%H:%MZ)-log 2>&1
0 5 * * * /home/balloon/backend/manage.py forecast_preprocess > /home/balloon/log/preprocess-$(date +%Y-%m-%dT%H:%MZ)-log 2>&1
0 2 * * * find /home/balloon/data -name '*.grib2' -mtime +2 -exec rm {} \;
